# stow/nfs/.config/systemd/user/home-zan-Music.automount
[Unit]
Description=Automount NFS Music

[Automount]
Where=/home/zan/Music
TimeoutIdleSec=10min

[Install]
WantedBy=multi-user.target
