# stow/nfs/.config/systemd/user/home-zan-Drive.automount
[Unit]
Description=Automount NFS Drive
After=network-online.target

[Automount]
Where=/home/zan/Drive
TimeoutIdleSec=10min

[Install]
WantedBy=multi-user.target
