# stow/nfs/.config/systemd/user/home-zan-Books.automount
[Unit]
Description=Automount NFS Books

[Automount]
Where=/home/zan/Books
TimeoutIdleSec=10min

[Install]
WantedBy=multi-user.target
