# stow/ghostty/.local/share/systemd/user/app-com.mitchellh.ghostty.service
[Unit]
Description=Ghostty
After=graphical-session.target
After=dbus.socket
Requires=dbus.socket

[Service]
Type=notify-reload
ReloadSignal=SIGUSR2
BusName=com.mitchellh.ghostty
ExecStart=/home/zan/bin/ghostty --gtk-single-instance=true --initial-window=false

[Install]
WantedBy=graphical-session.target
